# git-status.cron
*/10 * * * * /root/bsides-austin-2015/git-status.sh
